# deploy/crypto-trader.service
# systemd service สำหรับรันบอทค้างบน VPS (Oracle Cloud / Ubuntu)
# ติดตั้ง: ดูขั้นตอนใน DEPLOY_ORACLE.md
#
# ไฟล์นี้ให้บอท:
#   - รันครบ 5 คู่ แล้ววนซ้ำทุก 30 นาทีผ่าน deploy/run-bot-loop.sh
#   - รีสตาร์ทเองอัตโนมัติถ้า process ดับ (Restart=always)
#   - สตาร์ทเองหลัง VPS รีบูต (WantedBy=multi-user.target)

[Unit]
Description=Crypto Trader Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/crypto-trader
Environment=EXCHANGE_NAME=kraken
Environment=TIMEFRAME=1h
Environment=INTERVAL_SECONDS=1800
ExecStart=/usr/bin/env bash /home/ubuntu/crypto-trader/deploy/run-bot-loop.sh
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
